# ts-bsp/sources/meta-quickphone/recipes-core/systemd/files/usr-share-quick-phone-assets.mount
[Unit]
Description=User Data Directory for the Quickphone Project

[Mount]
What=/dev/mmcblk1p2
Where=/usr/share/quick-phone/assets
